# systemd/discord-updater.timer
[Unit]
Description=Run Discord updater daily and shortly after boot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target